# rootfs/etc/systemd/system/autopass.service
[Unit]
Description=Configure user account with autogenerated password
Conflicts=shutdown.target
After=systemd-user-sessions.service plymouth-quit-wait.service
After=rc-local.service rhel-readonly.service readonly-root.service
Before=getty.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/libexec/autopass.sh "$USER" "$PASSLEN" "$EXCLUDE_AMBIGUOUS"
Environment=USER=root
Environment=PASSLEN=8
EnvironmentFile=-/etc/sysconfig/autopass

[Install]
WantedBy=getty.target
